# autostart/rocket-notifications.desktop
[Desktop Entry]
Type=Application
Name=Rocket Notifications
Exec=/usr/bin/rocket-session-bin --notifications
X-GNOME-Autostart-enabled=true
