# autostart/rocket-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Rocket Wallpaper
Exec=/usr/bin/rocket-session-bin --wallpaper
X-GNOME-Autostart-enabled=true
